# icc/files/etc/hotplug.d/icc/00-test.sh
#!/bin/sh

#logger -s -p debug -t "icc" "hotplug: $INTERFACE $DEVICE $STATUS $1 $2"

