# icc/base-files/etc/uci-defaults/95_icc
#!/bin/sh
if [ ! -f /etc/config/icc ];then
	uci import icc < /dev/null
	uci -q batch <<-EOF >/dev/null
	delete icc.@icc[-1]
	add icc icc
	set icc.@icc[-1].enabled=0
	commit icc
	EOF
fi
exit 0
