# luci-app-icc/top-level-root-dir/etc/uci-defaults/95_icc
#!/bin/sh

#uci import icc < /dev/null

#uci -q batch <<-EOF >/dev/null
#	delete icc.@icc[-1]
#	add icc icc
#	set icc.@icc[-1].enabled=0
#	commit icc
#EOF
exit 0
